# templates/etc/profile.d/docker.sh.j2
# when we get rid of the docker socket we can turn the TLS bit on for everyone
# export DOCKER_TLS_VERIFY=1
